# tools/clawdbot-local-llm/launchagents/com.clawdbot.radeon-wrapper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.clawdbot.radeon-wrapper</string>

    <key>Comment</key>
    <string>Radeon LAN Tool Wrapper for ClawdBot (AMD RX 7900 XT)</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/bin/node</string>
      <string>/Users/mattiapetrucciani/CascadeProjects/clawdbot-config-ui/radeon-wrapper.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/mattiapetrucciani/CascadeProjects/clawdbot-config-ui</string>

    <key>StandardOutPath</key>
    <string>/Users/mattiapetrucciani/.clawdbot/logs/radeon-wrapper.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mattiapetrucciani/.clawdbot/logs/radeon-wrapper.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>/Users/mattiapetrucciani</string>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
  </dict>
</plist>
